# study-server.service
[Unit]
Description=Study Server - AI quiz/flashcard/summary generation
After=network.target

# Template systemd unit. Copy to /etc/systemd/system/ (or ~/.config/systemd/user/
# for a user unit), then replace the placeholders below with your own values:
#   <USER>       → the local user that should run the service
#   <INSTALL>    → absolute path to the cloned repository
# After editing, run:
#   systemctl daemon-reload
#   systemctl enable --now study-server

[Service]
Type=simple
User=<USER>
WorkingDirectory=<INSTALL>
ExecStart=/usr/bin/node <INSTALL>/dist/server.js
Restart=on-failure
RestartSec=5
Environment=STUDY_SERVER_PORT=3457
EnvironmentFile=<INSTALL>/.env

[Install]
WantedBy=multi-user.target
